# server-imgs.service
[Unit]
Description=Foo

[Service]
ExecStart=/bin/server-imgs
